# roles/build_firmware_sav/files/etc/uci-defaults/75-rootfs-repair
#!/bin/sh
set -x

# https://openwrt.org/docs/guide-user/advanced/expand_root?s[]=resize#instructions_manual
# OpenWrt 24.10+ - Resize root partition to 100% of disk

# Redirect debug output to log file
{
    [ -e /etc/rootfs-repair-ok ] && exit 0
    [ ! -e /etc/rootpt-resize-ok ] && exit 1

    lock -n /var/lock/root-resize

    ROOT_BLK="$(readlink -f /sys/dev/block/"$(awk -e '$9=="/dev/root"{print $3}' /proc/self/mountinfo)")"
    ROOT_PART="/dev/${ROOT_BLK##*/}"

    tune2fs -O^resize_inode "${ROOT_PART}"

    # Check and repair filesystem
    e2fsck -p "${ROOT_PART}"

    touch /etc/rootfs-repair-ok
} > /var/log/rootfs-repair 2>&1

reboot
