# roles/firmware_imagebuilder/templates/etc/uci-defaults/90-system.j2
#!/bin/sh

LOGFILE="/root/uci-defaults.log"

log() {
    echo "[$(date +'%Y-%m-%d %H:%M:%S')] [90-system] $*" >> "$LOGFILE"
}

log "======================================"
log "Starting 90-system configuration"
log "======================================"

# Detect if running in virtual environment
log "Detecting virtual environment"
grep -qEi "virtualbox|qemu" /sys/class/dmi/id/modalias 2>/dev/null
IS_VIRTUAL=$?
WAN_IFACE='eth3.100'
MGNT_SUBNET_ID='9'
if [ $IS_VIRTUAL -eq 0 ]; then
    WAN_IFACE='eth3'
    MGNT_SUBNET_ID='99'
    log "Virtual environment detected"
else
    log "Physical environment detected"
fi
log "WAN interface: $WAN_IFACE, MGNT subnet: $MGNT_SUBNET_ID"
uci set network.wan.device="$WAN_IFACE"
uci set network.wan6.device="$WAN_IFACE"
uci set network.mgnt.ipaddr="192.168.$MGNT_SUBNET_ID.1"
uci set network.mgnt.ip6hint="$MGNT_SUBNET_ID"
uci delete dhcp.mgnt.dhcp_option || true
uci set dhcp.mgnt.dhcp_option="6,192.168.$MGNT_SUBNET_ID.1"
uci commit
reload_config

#######################################
# change root password
#######################################

log "Changing root password"
passwd root <<EOP
{{ root_password }}
{{ root_password }}
EOP
log "Root password changed successfully"

# Replace ash by bash for root
log "Changing root shell from ash to bash"
sed -i "/^root:/ s|/bin/ash|/bin/bash|" /etc/passwd
log "Root shell changed to bash"

log "Setting permissions on /root/*.sh"
chmod u=rwx,g=,o= /root/*.sh
log "Permissions set successfully"

# Secure SSH access
log "Configuring SSH access"
mkdir -p /root/.ssh
chmod 700 /root/.ssh
echo {{ ssh_key }} >> /root/.ssh/authorized_keys
chmod 600 /root/.ssh/authorized_keys
log "SSH authorized_keys configured"
# Restart SSH to apply configuration
log "Restarting SSH service"
service sshd restart
log "SSH service restarted"

log "Enabling collectd sensors"
uci set luci_statistics.collectd_sensors.enable='1'
log "Collectd sensors enabled"

# Configure uHTTPd to listen only on MGNT interface
log "Configuring uHTTPd to listen only on MGNT interface"
MGNT_IPV4="192.168.$MGNT_SUBNET_ID.1"
MGNT_IPV6="{{ ula_prefix }}:$MGNT_SUBNET_ID::1"
uci del uhttpd.main.listen_http || true
uci del uhttpd.main.listen_https || true
uci add_list uhttpd.main.listen_http="$MGNT_IPV4:80"
uci add_list uhttpd.main.listen_http="[$MGNT_IPV6]:80"
uci add_list uhttpd.main.listen_https="$MGNT_IPV4:443"
uci add_list uhttpd.main.listen_https="[$MGNT_IPV6]:443"
log "uHTTPd configured: HTTP/HTTPS on $MGNT_IPV4 and $MGNT_IPV6"

log "Configuring SQM to use WAN interface"
uci set sqm.wan.interface="$WAN_IFACE"
log "SQM configured successfully"

log "Setting DDNS services enabled state"
uci set ddns.{{ ddns_service_name | replace('.', '_') }}_v4.enabled='$IS_VIRTUAL'
uci set ddns.{{ ddns_service_name | replace('.', '_') }}_v6.enabled='$IS_VIRTUAL'
log "DDNS enabled state set successfully"

log "Committing system configuration"
uci commit
reload_config
log "90-system configuration completed successfully"
log "======================================"
exit 0
