# roles/build_source_firmware/files/uci-defaults/90-bouygues-wan
#!/bin/sh
# Configure WAN on VLAN 100 for Bouygues FTTH (example) - idempotent

IF="eth0" # Adjust if your physical WAN differs; you can override by dropping another script later

# Only run once
[ -f /etc/config/network ] || exit 0

# Check if already configured
uci -q get network.wan && exit 0

uci batch <<EOF
set network.wan=interface
set network.wan.device="${IF}.100"
set network.wan.proto='dhcp'
set network.wan.metric='10'
EOF

# Ensure device with VLAN 100 exists
uci -q set network.wan_dev=device
uci -q set network.wan_dev.name="${IF}.100"
uci -q set network.wan_dev.type='8021q'
uci -q set network.wan_dev.ifname="${IF}"
uci -q set network.wan_dev.vid='100'

uci commit network

# DHCP client already covered by proto=dhcp
exit 0
