# roles/build_firmware/templates/etc/uci-defaults/20-kernel-cmdline.j2
#!/bin/sh
# 20-kernel-cmdline - append extra kernel cmdline params into /boot/grub/grub.cfg
# Generated by build_firmware role at image build. Runs once via uci-defaults.
# Idempotent: only appends if first param token not already present on a linux line.

EXTRAS="{{ firmware_kernel_extra_cmdline }}"
GRUBCFG="/boot/grub/grub.cfg"
[ -f "$GRUBCFG" ] || exit 0

# If all params already present on first linux line, exit.
first_linux_line=$(grep -m1 '^\s*linux' "$GRUBCFG" || true)
if [ -n "$first_linux_line" ]; then
  all_present=true
  for p in $EXTRAS; do
    printf '%s' "$first_linux_line" | grep -q "[[:space:]]$p" || all_present=false
  done
  [ "$all_present" = true ] && exit 0
fi

first_token=$(echo "$EXTRAS" | awk '{print $1}')
[ -n "$first_token" ] || exit 0

tmpfile=$(mktemp)
awk -v extras="$EXTRAS" -v firsttok="$first_token" '
  /^\s*linux/ {
    if (index($0, firsttok)==0) {
      print $0 " " extras
      next
    }
  }
  {print}
' "$GRUBCFG" > "$tmpfile" && mv "$tmpfile" "$GRUBCFG"

exit 0
