# roles/build_firmware/templates/etc/uci-defaults/91-system.j2
#!/bin/sh

#######################################
# change root password
#######################################

passwd root <<EOP
{{ root_password }}
{{ root_password }}
EOP

echo {{ ssh_key }} >> /etc/dropbear/authorized_keys
