# roles/build_image/files/etc/uci-defaults/90-admin
#!/bin/sh

# Enable DHCP client to wan
uci set network.wan.vendorid="neufbox_NB6VAC-directAccess"
uci commit network

exit 0
