# roles/build_firmware_sav/templates/etc/uci-defaults/81-mount-datadrive.j2
#!/bin/sh
. /lib/functions.sh

if [ ! -d {{ datadrive_mount_dir }} ] && lock -n /var/lock/mount-datadrive
then
    mkdir -p {{ datadrive_mount_dir }}

    # Find UUID of partition labeled "data"
    data_uuid=$(blkid -l -o value -s UUID -t LABEL=data)

    if [ -z "$data_uuid" ]; then
        echo "Error: Could not find partition labeled 'data'" >&2
        exit 1
    fi

    # Load fstab config and find mount section by UUID
    config_load fstab

    found_section=""

    _find_uuid_mount() {
        local section="$1"
        local uuid

        config_get uuid "$section" uuid

        # Check if UUID matches our data partition
        [ "$uuid" = "$data_uuid" ] && found_section="$section"
    }

    config_foreach _find_uuid_mount mount

    # Update if found
    if [ -n "$found_section" ]; then
        uci set fstab.$found_section.target='{{ datadrive_mount_dir }}'
        uci set fstab.$found_section.enabled='1'
        uci commit fstab
        reboot
    fi
fi
exit 0
