# uci-defaults/99-custom
#!/bin/sh
uci set network.wan.vendorid="neufbox_NB6VAC-directAccess"
uci commit network
exit 0